# hfs/conf/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/data/coze/run/supervisord.pid
childlogdir=/data/coze/logs

[unix_http_server]
file=/data/coze/run/supervisor.sock
chmod=0700
chown=cozeadmin:cozeadmin

[supervisorctl]
serverurl=unix:///data/coze/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:mariadb]
priority=10
directory=/app
command=/opt/coze-hfs/bin/run-mariadb.sh
user=user
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/data/coze/logs/mariadb.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/data/coze/logs/mariadb.err
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=2

[program:redis]
priority=20
command=/usr/bin/redis-server /opt/coze-hfs/conf/redis.conf
user=user
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/data/coze/logs/redis.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/data/coze/logs/redis.err
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=2

[program:nats]
priority=30
command=/usr/sbin/nats-server -c /opt/coze-hfs/conf/nats-server.conf
user=user
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/data/coze/logs/nats.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/data/coze/logs/nats.err
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=2

[program:minio]
priority=35
command=/opt/coze-hfs/bin/run-minio.sh
user=user
autostart=true
autorestart=false
startsecs=1
stdout_logfile=/data/coze/logs/minio.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/data/coze/logs/minio.err
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=2

[program:minio-init]
priority=36
command=/opt/coze-hfs/bin/init-minio.sh
user=user
autostart=true
autorestart=false
startsecs=0
stdout_logfile=/data/coze/logs/minio-init.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/data/coze/logs/minio-init.err
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=2

[program:etcd]
priority=40
command=/opt/coze-hfs/bin/run-etcd.sh
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/data/coze/logs/etcd.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/data/coze/logs/etcd.err
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=2

[program:elasticsearch]
priority=45
command=/opt/coze-hfs/bin/run-elasticsearch.sh
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/data/coze/logs/elasticsearch.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/data/coze/logs/elasticsearch.err
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=2

[program:milvus]
priority=50
directory=/milvus
command=/opt/coze-hfs/bin/run-milvus.sh
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/data/coze/logs/milvus.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/data/coze/logs/milvus.err
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=2

[program:coze-server]
priority=60
directory=/app
command=/opt/coze-hfs/bin/run-coze-server.sh
user=user
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/data/coze/logs/coze-server.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/data/coze/logs/coze-server.err
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=2

[program:ops-service]
priority=65
directory=/app
command=/opt/coze-hfs/bin/run-ops-service.sh
user=user
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/data/coze/logs/ops-service.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/data/coze/logs/ops-service.err
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=2

[program:admin-service]
priority=66
directory=/app
command=/opt/coze-hfs/bin/run-admin-service.sh
user=cozeadmin
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/data/coze/logs/admin-service.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/data/coze/logs/admin-service.err
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=2

[program:nginx]
priority=70
command=/usr/sbin/nginx -c /opt/coze-hfs/conf/nginx.conf -g 'daemon off;'
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/data/coze/logs/nginx.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/data/coze/logs/nginx.err
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=2
